# conf/duckdns.service
[Unit]
Description=Duck DNS Updater

[Service]
Type=simple
ExecStart=/home/userduckdns

[Install]
WantedBy=default.target
