# conf/duckdns.service
[Unit]
Description=Duck DNS Updater

[Service]
Type=simple
ExecStart=/home/user/duckdns

[Install]
WantedBy=default.target
